# debian/postinst
#!/bin/sh
set -e

# Enable and start user systemd services for all users
# This runs during package installation

if [ "$1" = "configure" ]; then
    # Reload systemd user daemon to pick up new unit files
    # We need to do this for each logged-in user
    for user_id in $(loginctl list-users --no-legend | awk '{print $1}'); do
        if [ -n "$user_id" ]; then
            # Get the username from the UID
            user_name=$(id -nu "$user_id" 2>/dev/null || true)
            if [ -n "$user_name" ]; then
                # Run systemctl commands as the user
                su - "$user_name" -c "systemctl --user daemon-reload" 2>/dev/null || true
                su - "$user_name" -c "systemctl --user enable battery-monitor.timer" 2>/dev/null || true
                su - "$user_name" -c "systemctl --user enable diskspace-monitor.timer" 2>/dev/null || true
                su - "$user_name" -c "systemctl --user start battery-monitor.timer" 2>/dev/null || true
                su - "$user_name" -c "systemctl --user start diskspace-monitor.timer" 2>/dev/null || true
            fi
        fi
    done
fi

#DEBHELPER#

exit 0
